# deployment/proxy-services/cyclopskernelwg.service
[Unit]
Description=Cyclops Wireguard interface
After=network.target network-online.target
Requires=network-online.target
 
[Service]
ExecStart=/deploy/kernelwg

[Install]
WantedBy=multi-user.target
